# scripts/fiber-kiosk.service
[Unit]
Description=Fiber Kiosk — LVGL touchscreen UI
After=network.target fiber-bridge.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/fiber-kiosk --fb /dev/fb0 --touch /dev/input/event0
Restart=always
RestartSec=3
Environment=DISPLAY=

[Install]
WantedBy=multi-user.target
